# .config/autostart/mako.desktop
[Desktop Entry]
Exec=/usr/bin/mako
Name=mako
Type=Application
