# .config/autostart/mpd.desktop
[Desktop Entry]
Exec=/usr/bin/mpd
Name=mpd
Type=Application
