# .config/autostart/pipewire.desktop
[Desktop Entry]
Exec=/usr/bin/pipewire
Name=pipewire
Type=Application
