# .config/autostart/pipewire-pulse.desktop
[Desktop Entry]
Exec=/usr/bin/pipewire-pulse
Name=pipewire-pulse
Type=Application
