# .config/autostart/syncthing.desktop
[Desktop Entry]
Exec=/usr/bin/syncthing -no-browser -no-restart
Name=syncthing
Type=Application
