# .config/autostart/wireplumber.desktop
[Desktop Entry]
Exec=/usr/bin/wireplumber
Name=wireplumber
Type=Application
