# .config/autostart/wlsunset.desktop
[Desktop Entry]
Exec=/usr/bin/wlsunset -l 44 -L 20
Name=wlsunset
Type=Application
